# recipes-bringup/rdkshell/files/wpeframework.service.in
[Unit]
Description=wpeframework
#Wants=multi-user.target
#After=multi-user.target  
After=local-fs.target 
After=memcr.service
DefaultDependencies=no
PartOf=multi-user.target

[Service]
PIDFile=/var/run/WPEFramework.pid
EnvironmentFile=-/etc/default/westeros-env
Environment="XDG_RUNTIME_DIR=/tmp"
Environment="LD_PRELOAD=libwesteros_gl.so.0.0.0"
Environment="WAYLAND_DISPLAY=wayland-0"
Environment="AAMP_ENABLE_OPT_OVERRIDE=TRUE"
Environment="GST_REGISTRY=/opt/.gstreamer/registry.bin"
Environment="GST_REGISTRY_UPDATE=no"
Environment="ENABLE_WEBKITBROWSER_PLUGIN_ACCESSIBILITY=1"
Environment="TTS_USE_THUNDER_CLIENT=1"
Environment="WESTEROS_SINK_LOW_MEM_MODE=1"
Environment="MALLOC_ARENA_MAX=1"
Environment="RDKSHELL_KEYMAP_FILE=/etc/rdkshell_keymapping.json"
Environment="RDKSHELL_FACTORY_APP_URL=http://localhost:50050/factoryui/index.html"
Environment="RDKSHELL_MAINBOARD_MANUFACTURING_FACTORY_APP_URL=http://localhost:50050/factoryui/index.html"
Environment="RDKSHELL_ASSEMBLY_FACTORY_APP_URL=http://localhost:50050/assemblyfactoryui/index.html"
Environment="RDKSHELL_ASSEMBLY_FACTORY_2_APP_URL=http://localhost:50050/assemblyfactoryui_2/index.html"
Environment="RDKSHELL_POWER_KEY_CODE=122"
Environment="RDKSHELL_WAIT_FOR_PERSISTENT_STORE=1"
Environment="RDKSHELL_BLOCK_RESIDENTAPP_FACTORYMODE=1"
Environment="RDKSHELL_FRAMERATE=60"
Environment="RDK_NO_ACTION_ON_POWER_KEY=1"
Environment="WESTEROS_GL_GRAPHICS_MAX_SIZE=1920x1080"
Environment="PLAYERSINKBIN_USE_WESTEROSSINK=1"
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/var/run/dbus/system_bus_socket
Environment="PLAYREADY_OP_DIGITAL_LEVEL=200"
Environment="WESTEROS_GL_USE_REFRESH_LOCK=1"
Environment="WESTEROS_SINK_USE_FREERUN=1"
Environment="WESTEROS_GL_DEBUG=2"
Environment="WESTEROS_SINK_USE_ESSRMGR=1"
Environment="WESTEROS_GL_MODE=3840x2160x50"
Environment="RDKSHELL_STARTUP_CONFIG=/etc/rdkshell_post_startup.conf"
Environment="GST_REGISTRY=/opt/.gstreamer/registry.bin"
#Environment="WESTEROS_GL_USE_GBM_MODIFIERS=1"
Environment="WESTEROS_GL_REFRESH_PRIORITY=R,70"
Environment="WESTEROS_GL_ALLOW_4K_ZOOM=0"
Environment="WESTEROS_GL_GLOBAL_ZOOM_ACTIVE=1"
Environment="WPE_WAL_AUTOCHECKPOINT=10"
Environment="OPEN_CDM_SERVER=/tmp/OCDM/ocdm"
Environment="GRPC_DEFAULT_SSL_ROOTS_FILE_PATH=/usr/share/ca-certificates/ca-certificates.crt"
ExecStartPre=/bin/sh -c '/lib/rdk/updatePreviousRebootInfo.sh; if [ -f /lib/rdk/logMilestone.sh ];then sh /lib/rdk/logMilestone.sh "WPE_FRAMEWORK_START"; fi;'
ExecStart=-/usr/bin/WPEFramework -b -c /etc/WPEFramework/config.json
ExecStopPost=/bin/sh -c '/usr/bin/killall -9 WPEProcess; touch /tmp/thunder_restarted'
KillSignal=SIGKILL
Restart=always

[Install]
WantedBy=ui-init.target
